# com.tidpl.scheduler.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tidpl.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/tidpl/run-if-due.sh</string>
    </array>

    <!-- Check every 6 hours -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/tidpl/data/scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/tidpl/data/scheduler.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
